# etc/rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#
# limit root cpu usage 
# echo 512 > /sys/kernel/uids/0/cpu_share
echo 2 > /proc/irq/1/smp_affinity
echo 1 > /proc/irq/12/smp_affinity
echo 2 > /proc/irq/14/smp_affinity
echo 2 > /proc/irq/18/smp_affinity
echo 2 > /proc/irq/19/smp_affinity
